# processmon/crontab
crontab -e
============================
* * * * * ~/bin/procmon start >> /var/log/dtmf/processmon.log 2>&1
============================
